# deploy/back.service
[Unit]
Description=FastAPI Backend Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/project/back
Environment="PATH=/var/www/project/back/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/var/www/project/back"
# Use wrapper script to ensure PYTHONPATH is set correctly
ExecStart=/var/www/project/back/start_backend.sh

Restart=always
RestartSec=5
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
